# voxtype/voxtype.service
# Launches the wrapper voxtype-with-dictionary.sh so the Whisper dictionary comes
# from ~/.config/voxtype/dictionary.txt.

[Unit]
Description=Voxtype voice-to-text daemon
Documentation=https://voxtype.io
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.config/voxtype/voxtype-with-dictionary.sh
Restart=on-failure
RestartSec=5

# Ensure we have access to the display
Environment=XDG_RUNTIME_DIR=%t

[Install]
WantedBy=graphical-session.target
